# packages/caramos-ota/etc/xdg/autostart/caramos-ota-notifier.desktop
[Desktop Entry]
Type=Application
Name=CaramOS Update Center
Comment=Check for and install CaramOS OTA updates
Exec=/usr/bin/caramos-ota-notifier --autostart
Icon=system-software-update
Terminal=false
NoDisplay=true
X-GNOME-Autostart-Phase=Applications
X-GNOME-Autostart-Delay=30
