# packages/caramos-ota/lib/systemd/system/caramos-ota-check.service
[Unit]
Description=CaramOS OTA check for updates
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/caramos-ota --check
PrivateTmp=true
ProtectHome=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
