# packages/caramos-ota/debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|deconfigure)
        # Stop and disable the timer on removal
        if command -v systemctl >/dev/null 2>&1; then
            systemctl stop caramos-ota-check.timer || true
            systemctl disable caramos-ota-check.timer || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
